# competition/servers/clojure-aleph.ini
[program:clojure-aleph]
command=java -jar cljwsdemo-0.1.0-SNAPSHOT-standalone.jar
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
