# competition/servers/erlang-cowboy-master.ini
[program:erlang-cowboy-master]
command=erl +P 134217727 -noshell -pa cowboy-master/deps/ranch/ebin cowboy-master/ebin +K true -s wsdemo_master
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
