# competition/servers/erlang-cowboy.ini
[program:erlang-cowboy]
command=erl -noshell -pa ../ebin ../deps/cowboy/ebin +K true -s wsdemo
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
