# competition/servers/erlang-cowboy.ini
[program:erlang-cowboy-0.6.1]
command=erl +P 134217727 -noshell -pa cowboy-0.6.1/deps/*/ebin cowboy-0.6.1/ebin +K true -s wsdemo
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
