# competition/servers/haskell-snap.ini
[program:haskell-snap]
command=./wsdemo-snap +RTS -N -A4M -qg1
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
