# competition/servers/java-netty4.ini
[program:java-netty4]
command=java -cp netty-all-4.0.0.Beta2-SNAPSHOT.jar:netty-example-4.0.0.Beta2-SNAPSHOT.jar io.netty.example.http.websocketx.autobahn.AutobahnServer 8000
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
