# competition/servers/node-ws-cluster.ini
[program:node-ws-cluster]
command=node wsdemo-ws-cluster.js
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
