# competition/servers/node-ws.ini
[program:node-ws]
command=node wsdemo-ws.js
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
