# competition/servers/perl-ev.ini
[program:perl-ev]
command=perl wsdemo.pl
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
