# competition/servers/pypy-tornado-N.ini
[program:pypy-tornado-N]
command=./pypy-1.9/bin/pypy wsdemo-tornado-N.py
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
