# competition/servers/pypy-tornado-1.ini
[program:pypy-tornado-1]
command=./pypy-1.9/bin/pypy wsdemo-tornado-1.py
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
