# competition/servers/pypy-twisted-N.ini
[program:pypy-twisted-N]
command=./pypy-1.9/bin/pypy wsdemo-twisted-N.py
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
