# competition/servers/python-twisted-1.ini
[program:python-twisted-1]
command=python wsdemo-twisted-1.py
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
