# competition/servers/python-gevent-websocket-N.ini
[program:python-gevent-websocket-N]
command=python wsdemo-gevent-websocket-N.py
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
