# competition/servers/python-tornado-1.ini
[program:python-tornado-1]
command=python wsdemo-tornado-1.py
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
