# competition/servers/python-twisted-N.ini
[program:python-twisted-N]
command=python wsdemo-twisted-N.py
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
