# competition/servers/python-tornado-N.ini
[program:python-tornado-N]
command=python wsdemo-tornado-N.py
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
