# competition/servers/ruby-em-websockets.ini
[program:ruby-em-websockets]
command=ruby wsdemo-em.rb
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
