# competition/servers/scala-play.ini
[program:scala-play]
command=bash ./wsdemo-playframework-2/wsdemo-1.0-SNAPSHOT/start -Dhttp.port=8000
autostart=false
autorestart=false
startretries=0
stopasgroup=true
killasgroup=true
